# overlay/etc/cron.daily/invoiceninja
#!/bin/bash

php /var/www/invoiceninja/artisan ninja:send-invoices
php /var/www/invoiceninja/artisan ninja:send-reminders
